# torcms_dde/pycsw_dist/pycsw_cfg/etc/harvest-all.cron
# run cronjob daily at 0400h
#
# info on setting up cron at http://www.unixgeeks.org/security/newbie/unix/cron-1.html

00 04 * * * /path/to/pycsw-admin.py refresh-harvested-records --config /path/to/default.cfg
